# ops/ai.openclaw.storykeeper-publisher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.storykeeper-publisher</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/djstd/.openclaw/storykeeper-library-runtime/ops/publish-queue.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>STORYKEEPER_LIBRARY</key>
    <string>/Users/djstd/.openclaw/storykeeper-library-runtime</string>
    <key>STORYKEEPER_QUEUE</key>
    <string>/Users/djstd/.openclaw/workspace-storyteller/data/publication-queue</string>
    <key>STORYKEEPER_PUBLISHED</key>
    <string>/Users/djstd/.openclaw/workspace-storyteller/data/publication-published</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>60</integer>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>/Users/djstd/.openclaw/workspace-storyteller/data/logs/publication.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/djstd/.openclaw/workspace-storyteller/data/logs/publication-error.log</string>
</dict>
</plist>
